# core/busybox-init/files/inittab
::sysinit:/usr/lib/init/rc.boot
::restart:/sbin/init
::shutdown:/usr/lib/init/rc.shutdown
::respawn:runsvdir -P /var/service 'log: ................................................................................................................................................................................................................................................................'
